# systemd-user/xfsettingsd.service
[Unit]
Description=XFCE XSettings Daemon

[Service]
ExecStart=/usr/bin/xfsettingsd
Restart=on-failure

[Install]
WantedBy=graphical-session.target
